# media-stack/quadlet/deluge.container
[Unit]
Description=Deluge (cliente torrent)
After=gluetun.service
Requires=gluetun.service

[Container]
Image=lscr.io/linuxserver/deluge:2.2.0
ContainerName=deluge

# Todo o tráfego (torrent e web UI) sai pela stack de rede do Gluetun —
# sem isso, o Deluge usaria a rede normal do host, sem VPN. Container
# que entra via "container:" não declara PublishPort= próprio (ver
# gluetun.container, onde 8112/6881 são publicados de verdade) nem
# Network=<nome>.network — os dois são mutuamente exclusivos com isso.
Network=container:gluetun

Volume=%h/.config/containers/volumes/media-stack/deluge/config:/config:Z
Volume=${MEDIA_DATA_DIR}:/data:Z
EnvironmentFile=%h/.config/containers/env/media-stack.env

HealthCmd=CMD-SHELL curl --silent --fail http://127.0.0.1:8112/ || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=30s
Notify=healthy

# homepage — exibição no dashboard (ver ../homepage/). Descoberta do
# tsdproxy fica no gluetun.container (dono da porta publicada), não
# aqui — ver comentário lá.
Label=homepage.group=Mídia
Label=homepage.name=Deluge
Label=homepage.icon=deluge.svg
Label=homepage.href=https://deluge.<seu-tailnet>.ts.net
Label=homepage.description="Cliente torrent (via Gluetun/VPN)"

[Service]
Restart=always
TimeoutStartSec=60

[Install]
WantedBy=default.target
