# media-stack/quadlet/dispatcharr-celery.container
[Unit]
Description=Dispatcharr Celery (worker de tarefas em segundo plano)
After=dispatcharr-postgres.service dispatcharr-redis.service dispatcharr.service
Requires=dispatcharr-postgres.service dispatcharr-redis.service dispatcharr.service

[Container]
Image=ghcr.io/dispatcharr/dispatcharr:latest
ContainerName=dispatcharr-celery
Network=dispatcharr-net.network
Entrypoint=/app/docker/entrypoint.celery.sh

# Mesmo /data do dispatcharr.container — ver nota sobre :z nele.
Volume=%h/.config/containers/volumes/media-stack/dispatcharr/data:/data:z
EnvironmentFile=%h/.config/containers/env/dispatcharr.env

# Vars específicas do worker — o Celery usa DISPATCHARR_PORT/DISPATCHARR_WEB_HOST
# só pra falar com o container do dispatcharr.container (DVR recording),
# não fazem parte do dispatcharr.env porque não são usadas pelo container web.
Environment=DISPATCHARR_PORT=9191
Environment=DJANGO_SETTINGS_MODULE=dispatcharr.settings
Environment=PYTHONUNBUFFERED=1

# Secret criado via: podman secret create dispatcharr-postgres-password ~/.config/containers/secrets/dispatcharr/postgres-password.txt
Secret=dispatcharr-postgres-password,type=env,target=POSTGRES_PASSWORD

# Sem healthcheck — mesmo padrão do docker-compose.yml oficial: um worker
# de fila em segundo plano não tem endpoint HTTP pra sondar. Restart=always
# cobre o processo cair; Requires= acima cobre a dependência de start.

[Service]
Restart=always
# Mesma imagem pesada do dispatcharr.container — ver comentário lá.
TimeoutStartSec=300

[Install]
WantedBy=default.target
