# media-stack/quadlet/dispatcharr.container
[Unit]
Description=Dispatcharr (IPTV, EPG, VOD) — modo AIO (tudo em um container)

[Container]
# Upstream não publica tag semver (só :latest, :dev e hash de commit) — ver
# seção Auto-update no README pra como isso é observado mesmo assim.
#
# Modo AIO (DISPATCHARR_ENV=aio), não o modular usado antes neste repo —
# Postgres e Redis rodam DENTRO deste mesmo container (via pg_ctl/uwsgi
# attach-daemon no entrypoint.sh da própria imagem), um único /data guarda
# tudo (banco em /data/db, chave Django em /data/jwt, config/gravações).
# Trocado de modular pra AIO a pedido do usuário — a imagem evoluiu desde
# o deploy modular original e passou a ter um entrypoint.sh único cobrindo
# os dois modos, tornando o AIO uma opção viável sem as 4 units/secret
# que o modular exigia.
Image=ghcr.io/dispatcharr/dispatcharr:latest
ContainerName=dispatcharr
PublishPort=9191:9191

Environment=DISPATCHARR_ENV=aio
Environment=DISPATCHARR_LOG_LEVEL=info

Volume=%h/.config/containers/volumes/media-stack/dispatcharr/data:/data:Z

HealthCmd=CMD-SHELL curl --silent --fail http://127.0.0.1:9191/ || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=30s
Notify=healthy

# tsdproxy discovery — expõe na tailnet (ver ../tsdproxy/)
Label=tsdproxy.enable=true
Label=tsdproxy.name=dispatcharr
Label=tsdproxy.port.web=443/https:9191/http

# homepage — exibição no dashboard (ver ../homepage/)
Label=homepage.group=Mídia
Label=homepage.name=Dispatcharr
Label=homepage.icon=dispatcharr.svg
Label=homepage.href=https://dispatcharr.<seu-tailnet>.ts.net
Label=homepage.description="Gerenciador de IPTV (streams, EPG, VOD)"

# Upstream só publica :latest — sem tag semver pra comparar, o WUD não
# teria sinal nenhum sem isso (ver README do wud, seção "Tags não-semver
# não são observadas"). wud.watch.digest compara o digest da imagem em
# vez de versão, então ainda avisa quando existe uma build nova.
Label=wud.watch=true
Label=wud.watch.digest=true

[Service]
Restart=always
# Imagem pesada e primeiro boot em modo AIO faz MAIS coisa que o modular
# (initdb do Postgres embutido + migração + Redis + Celery + nginx +
# uwsgi, tudo em sequência) — 300s testado como suficiente pro modular,
# mantido aqui com folga de propósito.
TimeoutStartSec=300

[Install]
WantedBy=default.target
