# homepage/quadlet/homepage.container
[Unit]
Description=Homepage - dashboard de serviços
Requires=podman.socket
After=podman.socket

[Container]
Image=ghcr.io/gethomepage/homepage:latest
AutoUpdate=registry
ContainerName=homepage
PublishPort=3000:3000

# Socket compatível com a API do Docker, exposto pelo podman.socket —
# leitura de labels apenas (ro). Mesmo mecanismo já usado pelo tsdproxy.
Volume=%t/podman/podman.sock:/var/run/docker.sock:ro,z
Volume=%h/.config/containers/volumes/homepage/config:/app/config:Z
# Ícones customizados (não presentes no dashboard-icons nem como URL) —
# referenciar em Label=homepage.icon=/icons/<arquivo> pra usar.
Volume=%h/.config/containers/volumes/homepage/icons:/app/public/icons:Z
EnvironmentFile=%h/.config/containers/env/homepage.env

HealthCmd=CMD-SHELL wget -q --spider http://127.0.0.1:3000/ || exit 1
HealthInterval=10s
HealthTimeout=5s
HealthRetries=5
HealthStartPeriod=10s
Notify=healthy

# tsdproxy discovery — expõe o dashboard na tailnet (ver ../tsdproxy/)
Label=tsdproxy.enable=true
Label=tsdproxy.name=homepage
Label=tsdproxy.port.web=443/https:3000/http

[Service]
Restart=always

[Install]
WantedBy=default.target
