# radicale/birthday-sync/radicale-birthday-sync.service
[Unit]
Description=Regenera o calendário de aniversários do Radicale a partir dos contatos
Requires=radicale.service
After=radicale.service

[Service]
Type=oneshot
ExecStart=/usr/bin/podman exec --user radicale -w /data/collections radicale sh -c 'find collection-root -mindepth 3 -name "*.vcf" ! -path "*/.Radicale.cache/*" | python3 /data/birthday-calendar/create_birthday_calendar.py'
