# baikal/baikal.container
[Unit]
Description=Baikal (CalDAV/CardDAV)

[Container]
Image=docker.io/ckulka/baikal:0.10.1-nginx-php8.2
ContainerName=baikal
PublishPort=8084:80
Volume=%h/.config/containers/volumes/baikal/config:/var/www/baikal/config:Z
Volume=%h/.config/containers/volumes/baikal/data:/var/www/baikal/Specific:Z
EnvironmentFile=%h/.config/containers/env/baikal.env

HealthCmd=CMD-SHELL curl --silent --fail http://127.0.0.1:80/ || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=15s
Notify=healthy

# tsdproxy discovery — expõe na tailnet (ver ../tsdproxy/)
Label=tsdproxy.enable=true
Label=tsdproxy.name=baikal
Label=tsdproxy.port.web=443/https:80/http

# homepage — exibição no dashboard (ver ../homepage/)
Label=homepage.group=Self-Hosted
Label=homepage.name=Baikal
Label=homepage.icon=baikal.svg
Label=homepage.href=https://baikal.<seu-tailnet>.ts.net
Label=homepage.description="Servidor CalDAV/CardDAV self-hosted"

# A imagem também publica tags sem sufixo de variante (ex.: 0.10.1, sem
# "-nginx-phpX.Y") — o WUD trata isso como semver "sem prerelease", que
# semver considera "maior" que a tag com sufixo, gerando falso positivo
# de atualização. Restringe a candidatos da mesma família nginx+PHP.
# Regex sem \d/\. de propósito: o parser do Quadlet não aceita barra
# invertida em Label= ("unsupported escape char" no journal, dropa a
# linha inteira em silêncio, sem erro visível em `systemctl cat`).
Label=wud.watch=true
Label=wud.tag.include=^[0-9]+.[0-9]+.[0-9]+-nginx-php[0-9.]+$

[Service]
Restart=always

[Install]
WantedBy=default.target
