# any-sync-bundle/watchdog/any-sync-bundle-watchdog.timer
[Unit]
Description=Periodic any-sync-bundle liveness check

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
