# any-sync-bundle/quadlet/any-sync-mongo-legacy.container
# MongoDB 4.4 for CPUs without AVX instruction support (required by MongoDB 5.0+).
# More info: https://github.com/grishy/any-sync-bundle/pull/39
# Official recommendation: https://github.com/anyproto/any-sync-dockercompose/wiki/Troubleshooting-&-FAQ
#
# Use this instead of any-sync-mongo.container (not both — same ContainerName/NetworkAlias).
# To enable: don't copy any-sync-mongo.container into the systemd unit dir,
# copy this one instead renamed to any-sync-mongo.container.

[Unit]
Description=Any-Sync-Bundle MongoDB (legacy, no-AVX)

[Container]
Image=docker.io/mongo:4.4
ContainerName=any-sync-bundle-mongo
Network=any-sync-bundle-net.network
NetworkAlias=mongo
Exec=--replSet rs0 --port 27017
Volume=%h/.config/containers/volumes/any-sync-bundle/mongo:/data/db:Z
HealthCmd=CMD-SHELL RESULT=$$(mongo --quiet --host 127.0.0.1:27017 --eval 'try { if (rs.status().ok === 1) 1; else { rs.initiate({_id:"rs0", members:[{_id:0, host:"mongo:27017"}]}); rs.status().ok; } } catch (e) { rs.initiate({_id:"rs0", members:[{_id:0, host:"mongo:27017"}]}); rs.status().ok; } }') ; test "$$RESULT" = "1"
HealthInterval=5s
HealthTimeout=5s
HealthRetries=12
HealthStartPeriod=10s
Notify=healthy

# Restringe a candidatos dentro da própria linha 4.4.x — a 5.0+ exige
# AVX (motivo de usar esta imagem legacy), então uma tag 5.x/6.x não
# seria aplicável aqui de qualquer forma.
Label=wud.watch=true
Label=wud.tag.include=^4.4.[0-9]+$

[Service]
Restart=always
TimeoutStartSec=120

[Install]
WantedBy=default.target
