# any-sync-bundle/quadlet/any-sync-bundle.container
[Unit]
Description=Any-Sync-Bundle server
After=any-sync-mongo.service any-sync-bundle-redis.service
Requires=any-sync-mongo.service any-sync-bundle-redis.service

[Container]
# Tag explícita (recomendação do próprio projeto) — sem auto-update: esta
# imagem não tem shell/utilitários (é minimal de propósito), então não dá
# pra configurar HealthCmd nela, e sem HealthCmd não existe rollback
# automático do Podman. Bump manual: trocar a versão aqui, daemon-reload,
# restart. Ver Auto-update no README antes de mudar isso.
Image=ghcr.io/grishy/any-sync-bundle:1.4.3-2026-04-21-minimal
ContainerName=any-sync-bundle
Network=any-sync-bundle-net.network
PublishPort=33010:33010
PublishPort=33020:33020/udp
Volume=%h/.config/containers/volumes/any-sync-bundle/bundle:/data:Z
EnvironmentFile=%h/.config/containers/env/any-sync-bundle.env

# tsdproxy discovery — expõe o bundle na tailnet (TCP + QUIC/UDP)
Label=tsdproxy.enable=true
Label=tsdproxy.name=any-sync
Label=tsdproxy.port.sync=33010/tcp:33010/tcp
Label=tsdproxy.port.quic=33020/udp:33020/udp

# homepage — exibição no dashboard (ver ../homepage/)
Label=homepage.group=Self-Hosted
Label=homepage.name=Any-Sync-Bundle
Label=homepage.icon=anytype.svg
Label=homepage.description="Backend self-hosted do Anytype"

Label=wud.watch=true

[Service]
Restart=always

[Install]
WantedBy=default.target
