# owntracks/owntracks-mosquitto.container
[Unit]
Description=Mosquitto (broker MQTT do OwnTracks)

[Container]
Image=docker.io/library/eclipse-mosquitto:2.1.2-alpine
ContainerName=mosquitto
Network=owntracks-net.network
NetworkAlias=mosquitto
PublishPort=1883:1883
# WebSockets — ver mosquitto.conf
PublishPort=9001:9001

# mosquitto.conf fica no diretório de config — a imagem já usa
# /mosquitto/config/mosquitto.conf automaticamente se o arquivo existir,
# sem precisar de Exec= próprio.
Volume=%h/.config/containers/volumes/owntracks/mosquitto/config:/mosquitto/config:Z
Volume=%h/.config/containers/volumes/owntracks/mosquitto/data:/mosquitto/data:Z

# Secret criado via: podman secret create owntracks-mosquitto-passwd ~/.config/containers/secrets/owntracks/mosquitto-passwd
# uid/gid=1883: uid interno do usuário "mosquitto" que o processo usa
# depois de subir. mode=0400: o Mosquitto reclama em dois avisos
# separados testados na prática — dono errado (default root:root do
# Podman) E permissão "mundo-legível" (0444 default também); sem os
# dois ajustes juntos ("Future versions will refuse to load this file").
Secret=owntracks-mosquitto-passwd,target=/mosquitto/secrets/passwd,uid=1883,gid=1883,mode=0400

HealthCmd=CMD-SHELL nc -z 127.0.0.1 1883 || exit 1
HealthInterval=15s
HealthTimeout=5s
HealthRetries=5
HealthStartPeriod=10s
Notify=healthy

# SEM tsdproxy aqui, de propósito — testado na prática: nó TCP puro
# (tsdproxy.port.*) trava permanentemente em "Starting"/NeedsLogin nesta
# versão do tsdproxy (2.3.4), nunca chega a "Running". O any-sync-bundle
# (único outro proxy TCP puro deste repositório) mostra o mesmo sintoma
# nos logs (health check falhando em loop) — não é specific deste app,
# parece limitação/bug do tsdproxy com proxy não-HTTP. PublishPort=1883
# acima já expõe a porta em todas as interfaces do host, incluindo a da
# tailnet — usar o hostname do PRÓPRIO HOST (não um nó tsdproxy) pra
# alcançar o broker (ver README, seção "Configurando o app").

[Service]
Restart=always
TimeoutStartSec=60

[Install]
WantedBy=default.target
