# nginx/nginx.container
[Unit]
Description=nginx (servidor de arquivos estáticos)

[Container]
Image=docker.io/library/nginx:1.30.4-alpine
ContainerName=nginx
Network=nginx-net.network
PublishPort=8086:80

# Conteúdo estático — editar os arquivos direto aqui, o nginx só lê (:ro)
Volume=%h/.config/containers/volumes/nginx/html:/usr/share/nginx/html:ro,Z

# server blocks customizados (opcional) — vazio por padrão, usa o
# nginx.conf/default.conf embutidos da imagem se nada for colocado aqui
Volume=%h/.config/containers/volumes/nginx/conf.d:/etc/nginx/conf.d:ro,Z

HealthCmd=CMD-SHELL wget -q --spider http://127.0.0.1:80/ || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=10s
Notify=healthy

# tsdproxy discovery — expõe na tailnet (ver ../tsdproxy/)
Label=tsdproxy.enable=true
Label=tsdproxy.name=nginx
Label=tsdproxy.port.web=443/https:80/http

# homepage — exibição no dashboard (ver ../homepage/)
Label=homepage.group=Self-Hosted
Label=homepage.name=nginx
Label=homepage.icon=nginx.svg
Label=homepage.href=https://nginx.<seu-tailnet>.ts.net
Label=homepage.description="Servidor de arquivos estáticos"

Label=wud.watch=true
# nginx publica duas branches em paralelo, minor ímpar = mainline
# (dev), par = stable — ambas com tags -alpine numeradas, então só
# restringir ao sufixo "-alpine" não basta (1.31.x-alpine também bate).
# Preso à minor 1.30 (stable atual) — trocar aqui quando o próprio
# projeto anunciar a próxima stable, não seguir mainline automático.
Label=wud.tag.include=^1.30.[0-9]+-alpine$

[Service]
Restart=always
TimeoutStartSec=60

[Install]
WantedBy=default.target
