# linkwarden/pgdump/linkwarden-pgdump.timer
[Unit]
Description=Roda o dump do Postgres do Linkwarden diariamente, antes do job do Zerobyte

[Timer]
OnCalendar=*-*-* 02:50:00
Persistent=true

[Install]
WantedBy=timers.target
