# any-sync-bundle/backup-pause/any-sync-bundle-backup-resume.service
[Unit]
Description=Religa any-sync-bundle depois do job de backup do bundle no Zerobyte

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user start any-sync-bundle.service
