# any-sync-bundle/backup-pause/any-sync-bundle-backup-pause.timer
[Unit]
Description=Agenda a parada do stack any-sync-bundle antes do job de backup no Zerobyte

[Timer]
OnCalendar=*-*-* 22:00:00
# Persistent=false de propósito: se o host esteve desligado nesse horário,
# não faz sentido "recuperar" uma parada atrasada — melhor pular e esperar
# a próxima janela agendada. Contraste com o timer de resume, que usa
# Persistent=true (ver any-sync-bundle-backup-resume.timer).
Persistent=false

[Install]
WantedBy=timers.target
