# any-sync-bundle/backup-pause/any-sync-bundle-backup-pause.service
[Unit]
Description=Para o stack any-sync-bundle (bundle+mongo+redis) antes do job de backup no Zerobyte

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user stop any-sync-bundle.service any-sync-mongo.service any-sync-bundle-redis.service
