# any-sync-bundle/backup-pause/any-sync-bundle-backup-resume.service
[Unit]
Description=Religa o stack any-sync-bundle (bundle+mongo+redis) depois do job de backup no Zerobyte

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user start any-sync-mongo.service any-sync-bundle-redis.service any-sync-bundle.service
