# tsdproxy/quadlet/tsdproxy.container
[Unit]
Description=tsdproxy - Tailscale reverse proxy manager
Requires=podman.socket
After=podman.socket

[Container]
Image=docker.io/almeidapaulopt/tsdproxy:2
ContainerName=tsdproxy

PublishPort=8080:8080

# Socket compatível com a API do Docker, exposto pelo podman.socket
Volume=%t/podman/podman.sock:/var/run/docker.sock:z
Volume=%h/.config/containers/volumes/tsdproxy/data:/data:Z
Volume=%h/.config/containers/volumes/tsdproxy/config:/config:Z

AddHost=host.docker.internal:host-gateway

# Secret criado via: podman secret create authkey ~/.config/containers/secrets/tsdproxy/authkey.txt
Secret=authkey

# Auto-descoberta do próprio dashboard (ver docs do tsdproxy)
Label=tsdproxy.enable=true
Label=tsdproxy.name=dash
Label=tsdproxy.dash.visible=false

# homepage — exibição no dashboard (ver ../homepage/)
Label=homepage.group=Infra
Label=homepage.name=tsdproxy
Label=homepage.icon=si-tailscale
Label=homepage.href=https://dash.<seu-tailnet>.ts.net
Label=homepage.description="Publica containers na tailnet automaticamente"

[Service]
Restart=always

[Install]
WantedBy=default.target
