# tsdproxy/quadlet/tsdproxy.container
[Unit]
Description=tsdproxy - Tailscale reverse proxy manager
Requires=podman.socket
After=podman.socket

[Container]
Image=docker.io/almeidapaulopt/tsdproxy:2
ContainerName=tsdproxy

PublishPort=8080:8080

# Socket compatível com a API do Docker, exposto pelo podman.socket
Volume=%t/podman/podman.sock:/var/run/docker.sock:z
Volume=%h/.config/containers/volumes/tsdproxy/data:/data:Z
Volume=%h/.config/containers/volumes/tsdproxy/config:/config:Z

AddHost=host.docker.internal:host-gateway

# Secret criado via: podman secret create authkey ~/.config/containers/secrets/tsdproxy/authkey.txt
Secret=authkey

# Auto-descoberta do próprio dashboard (ver docs do tsdproxy)
Label=tsdproxy.enable=true
Label=tsdproxy.name=dash
Label=tsdproxy.dash.visible=false

# homepage — exibição no dashboard (ver ../homepage/)
Label=homepage.group=Infra
Label=homepage.name=tsdproxy
Label=homepage.icon=tailscale.svg
Label=homepage.href=https://dash.<seu-tailnet>.ts.net
Label=homepage.description="Publica containers na tailnet automaticamente"

# A tag "3" no Docker Hub hoje resolve pra 3.0.0-beta.3 — sem release
# estável da v3 ainda, então "atualização disponível" não é acionável
# (não é bug de parsing tipo baikal/HA, a "3" existe de verdade, só que
# é beta). Restringe a candidatos com formato de release completo
# (x.y.z, sem sufixo -beta/-alpha) — exclui as tags de major sozinho
# (2, 3) e qualquer pre-release. Regex sem \d/\. de propósito, ver regra
# 18 do README raiz.
Label=wud.watch=true
Label=wud.tag.include=^[0-9]+.[0-9]+.[0-9]+$

[Service]
Restart=always

[Install]
WantedBy=default.target
